# scripts/systemd/safir-boot-splash.service
[Unit]
Description=SAFIR OLED Boot-Splash (zeigt System-Start auf OLED)
# Nach dem i2c-Kernelmodul und basis-filesystem läuft, aber so früh wie möglich
After=systemd-modules-load.service local-fs.target
DefaultDependencies=yes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/safir-oled-msg "UBUNTU" "startet..."
# Nach multi-user.target ist das OS ready, SAFIR wird danach geladen.
# Kein Fragezeichen, das ist eine Status-Anzeige kein Dialog.
ExecStartPost=/bin/sh -c 'sleep 2 && /usr/local/bin/safir-oled-msg "UBUNTU" "bereit"'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
